# packages/wbm-testbed/files/etc/uci-defaults/firstboot.sh
#!/bin/sh
echo "Configuring WBM testbed for first boot..."
lua /usr/bin/wbm-config
